# 11_OS + VSCode/1_Ubuntu/1_Deploy/3_Django-VPS-Deploy/Files_support/dev1.kopi34.ru.gunicorn.socket
[Unit]
Description=telebot7.kopi34.ru.gunicorn socket

[Socket]
ListenStream=/run/telebot7.kopi34.ru.gunicorn.sock

[Install]
WantedBy=sockets.target
